# render/worker/gaturro-render.timer
[Unit]
Description=EspermaRender cada 10 min
[Timer]
OnBootSec=5min
OnUnitActiveSec=10min
[Install]
WantedBy=timers.target
