# render/worker/gaturro-render.service
[Unit]
Description=EspermaRender worker (replay -> mp4 -> sotano)
[Service]
Type=oneshot
ExecStart=%h/.local/bin/gaturro-render-worker
